# systemd/systemd-age-creds.service
[Unit]
Description=systemd age credentials service
Documentation=https://github.com/abyrne55/systemd-age-creds
Requires=systemd-age-creds.socket
After=systemd-age-creds.socket
ConditionPathExists=/etc/age/identity.txt
ConditionDirectoryNotEmpty=/etc/age/credentials

[Service]
Type=simple
ExecStart=/usr/bin/systemd-age-creds
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/run
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true

# Default paths (override with systemctl edit)
Environment=AGE_IDENTITY=/etc/age/identity.txt
Environment=AGE_DIR=/etc/age/credentials

[Install]
WantedBy=multi-user.target
